# files/systemd/user/xsettingsd.service.d/10-custom.conf
[Unit]
After=xwayland-satellite.service
PartOf=xwayland-satellite.service

